# docker/backend/server.conf
[program:server]
command=php artisan octane:start --server=frankenphp --caddyfile=/usr/src/Caddyfile
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/server.log